# jukebox.socket
[Unit]
Description=Jukebox gunicorn socket

[Socket]
ListenStream=/run/gunicorn.sock

[Install]
WantedBy=sockets.target